# config/harbouros-flux-ipwatch.service
[Unit]
Description=Flux node external IP change detection and auto-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harbouros-flux-ipwatch.sh
