# config/harbouros-self-update.service
[Unit]
Description=Update HarbourOS from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harbouros-self-update.sh
