# config/harbouros-os-update.timer
[Unit]
Description=Weekly OS package updates

[Timer]
OnCalendar=Wed *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
